# systemd/emqx-monitor.service
[Unit]
Description=EMQX Connection Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/emqx/emqx-monitor
ExecStart=/opt/emqx/emqx-monitor/scripts/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin"

# 安全加固
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
